# systemd/bluesnap.service
[Unit]
Description=Bluesnap Bluetooth Snapcast Bridge
After=network-online.target bluetooth.target
Wants=network-online.target bluetooth.target
StartLimitInterval=30

[Service]
Type=simple
WorkingDirectory={{REPO_PATH}}
ExecStart={{REPO_PATH}}/scripts/bluesnap-boot.sh
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target

